# deploy/consorcio-sorteado.service
[Unit]
Description=Consórcio Sorteado — API de Automação Comercial
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=APP_USER
Group=APP_USER
WorkingDirectory=/opt/consorcio-sorteado
EnvironmentFile=/opt/consorcio-sorteado/.env
ExecStart=/opt/consorcio-sorteado/.venv/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --log-level info \
    --access-log
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=consorcio-sorteado

[Install]
WantedBy=multi-user.target
